# systemd/.config/systemd/user/omarchy-wallpaper.service
[Unit]
Description=Rotate Omarchy wallpaper via awww
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/tm-next-wallpaper-awww

